# config/systemd/user/battery-notifier.timer
[Unit]
Description=Battery notifier scheduler

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
